# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-firstboot.service.d/10-params.conf
[Service]
ExecStart=
ExecStart=systemd-firstboot --prompt-hostname --prompt-locale --prompt-keymap-auto --prompt-timezone --prompt-root-password --mute-console=yes
