# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-homed-firstboot.service.d/10-params.conf
[Service]
ExecStart=
ExecStart=homectl firstboot --prompt-new-user --storage=subvolume --prompt-shell=no --prompt-groups=no --mute-console=yes --uid=1000 --member-of=adm,systemd-journal,sudo,docker,incus-admin
