# mkosi.extra.common/usr/lib/systemd/system/opt-incus.mount
[Unit]
Before=incus.socket

[Mount]
What=/usr/share/factory/opt/incus
Where=/opt/incus
Options=bind
Type=none

[Install]
WantedBy=incus.socket
