# mkosi.extra.common/usr/lib/systemd/system/systemd-firstboot.service.d/10-params.conf
[Service]
ExecStart=
ExecStart=systemd-firstboot --prompt-hostname --prompt-locale --prompt-keymap --prompt-timezone --prompt-root-password --mute-console=yes
